# wakeup-check-post.service
# /etc/systemd/system/wakeup-check-post.service
[Unit]
Description=Wakeup Check Post Suspend
After=sleep.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'sleep 1; /usr/local/bin/wakeup-check.sh post'

[Install]
WantedBy=sleep.target
