# wakeup-check-post.service
[Unit]
Description=Wakeup Check Post Resume
After=suspend.target
After=hibernate.target
After=hybrid-sleep.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wakeup-check.sh post
TimeoutStartSec=300

[Install]
WantedBy=sleep.target
