# wakeup-check-pre.service
[Unit]
Description=Wakeup Check (Pre)
Before=sleep.target
Requires=sleep.target
#After=wakeup-check-post.service
#Conflicts=wakeup-check-post.service

[Service]
Type=oneshot
ExecStart=bash /usr/local/bin/wakeup-check.sh pre
RemainAfterExit=no

[Install]
WantedBy=sleep.target
