# wakeup-check-post.service
[Unit]
Description=Wakeup Check Post Resume
After=suspend.target
Conflicts=wakeup-check-pre.service
After=wakeup-check-pre.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wakeup-check.sh post
TimeoutStartSec=300

[Install]
WantedBy=suspend.target
WantedBy=post-suspend.target
