# wakeup-check-post.service
[Unit]
Description=Wakeup Check (Post)
After=suspend.target
#After=wakeup-check-pre.service
#Conflicts=wakeup-check-pre.service

[Service]
Type=oneshot
ExecStart=bash /usr/local/bin/wakeup-check.sh post
RemainAfterExit=no

[Install]
WantedBy=suspend.target
