# wakeup-check-post.service
[Unit]
Description=Wakeup Check - Post Resume
DefaultDependencies=no
After=suspend.target
Wants=suspend.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/wakeup-check.sh post

[Install]
WantedBy=suspend.target
